# afm-api-self_QML@.socket
#self generated by Rolando
[Unit]
Description=Provides websocket api self_QML for user %i
Requires=afm-user-setup@%i.service
After=afm-user-setup@%i.service
DefaultDependencies=no
[Socket]
SmackLabel=*
SmackLabelIPIn=System
SmackLabelIPOut=System
ListenStream=/run/user/%i/apis/ws/HVAC
FileDescriptorName=HVAC
Service=afm-service-agl-service-self_QML--1.0--main@%i.service
